# tools/weather-machine/infra/systemd/feed-sync.timer
[Unit]
Description=Refresh feed.json from S3 every minute
Documentation=https://github.com/inquirer/weather-machine

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
