# tools/weather-machine/infra/systemd/feed-sync.service
[Unit]
Description=Sync feed.json from S3 into the Weather Machine web docroot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=weatherweb
Group=weatherweb
EnvironmentFile=/etc/weather-machine/web.env
ExecStart=/opt/weather-machine/feed-sync.sh

# Hardening - may write only the docroot.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/weather-machine/site
ProtectHome=true
PrivateTmp=true
